# extras/systemd-timer/ncddns.service
[Unit]
Description=Updates DynDNS entry

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/ncddns/ncddns.py
